# focus-ledger.service
[Unit]
Description=Focus Ledger Dash App
After=network.target

[Service]
User=mail_
WorkingDirectory=/home/mail_/FocusLedger
ExecStart=/home/mail_/FocusLedger/.venv/bin/python -m focusledger.app
Environment=PYTHONUNBUFFERED=1
Restart=always

[Install]
WantedBy=multi-user.target